# iptables_rules/iptables.cron
* 12 * * *  cat  /var/log/messages |awk '{print $11}'|sed 's/SRC=//'|sort -n|uniq|xargs -I[] iptables -A INPUT  -p tcp -s [] -j DROP
2 12 * * * iptables-save
